# scripts/rc-proxy-tunnel.service
# systemd service for the WeCom SOCKS5 proxy tunnel.
#
# Install:
#   1. Copy this file to /etc/systemd/system/rc-proxy-tunnel.service
#   2. Edit ExecStart with your actual values (Elastic IP, port, key path)
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now rc-proxy-tunnel
#
# Check status:
#   sudo systemctl status rc-proxy-tunnel
#   journalctl -u rc-proxy-tunnel -f

[Unit]
Description=Remote Control — WeCom SOCKS5 proxy tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# --- EDIT THESE VALUES ---
ExecStart=/opt/remote_control/scripts/start-proxy-tunnel.sh ELASTIC_IP 1080 /home/deploy/.ssh/rc-proxy-key.pem
# -------------------------
Restart=always
RestartSec=5
User=deploy

[Install]
WantedBy=multi-user.target
